# rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/usr/bin/with-contenv bashio
# Home Assistant Community Add-on: n8n
# Initializes n8n environment variables and manages .env file

chmod +x "$0"

bashio::log.info "Initializing n8n configuration..."

N8N_DATA_DIR="/home/node/.n8n"
ENV_FILE="${N8N_DATA_DIR}/.env"
mkdir -p "${N8N_DATA_DIR}"
chown -R root:root "${N8N_DATA_DIR}"
chmod -R 700 "${N8N_DATA_DIR}"

# Fonction pour mettre à jour ou ajouter une variable dans .env
update_env_file() {
  local key="$1"
  local value="$2"
  if grep -q "^${key}=" "${ENV_FILE}"; then
    # Si la clé existe, la mettre à jour
    sed -i "s|^${key}=.*|${key}=${value}|" "${ENV_FILE}"
  else
    # Sinon, l'ajouter
    echo "${key}=${value}" >> "${ENV_FILE}"
  fi
}

# Si .env n'existe pas, le créer avec des valeurs par défaut
if [ ! -f "${ENV_FILE}" ]; then
  bashio::log.info "Création de ${ENV_FILE} avec des valeurs par défaut..."
  cat <<EOV > "${ENV_FILE}"
# Fichier de configuration n8n - éditez ces valeurs selon vos besoins
# Voir https://docs.n8n.io/hosting/configuration/environment-variables/
N8N_LOG_LEVEL=info
N8N_COMMUNITY_PACKAGES_ENABLED=false
N8N_PORT=5678
N8N_LISTEN_ADDRESS=0.0.0.0
N8N_USER_FOLDER=/home/node/.n8n
EOV
fi

# Lire les options depuis config.yaml
webhook_url=$(bashio::config 'webhook_url')
# Extraire protocol et host depuis webhook_url
protocol=$(echo "$webhook_url" | grep -Eo '^https?')
host=$(echo "$webhook_url" | sed -E 's|^https?://||' | sed -E 's|[:/].*||')

# Mettre à jour .env avec les valeurs de l'onglet Configuration
update_env_file "N8N_WEBHOOK_URL" "${webhook_url}"
update_env_file "N8N_PROTOCOL" "${protocol}"
update_env_file "N8N_HOST" "${host}"

# Lire et appliquer les variables personnalisées depuis la liste env
if bashio::config.has_value 'env'; then
  bashio::config 'env' | while read -r env_line; do
    if [[ -n "$env_line" ]]; then
      key=$(echo "$env_line" | cut -d'=' -f1)
      value=$(echo "$env_line" | cut -d'=' -f2-)
      update_env_file "$key" "$value"
    fi
  done
fi

bashio::log.info "Fichier .env mis à jour avec les valeurs de configuration :"
cat "${ENV_FILE}" | while read -r line; do bashio::log.info "$line"; done

exit 0